# deploy/ihostmc-website.service
[Unit]
Description=iHostMC Website (Next.js)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/iHostMC/website
Environment=NODE_ENV=production
Environment=PATH=/home/ubuntu/.nvm/versions/node/v22.22.0/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
EnvironmentFile=/opt/iHostMC/website/.env
ExecStart=/home/ubuntu/.nvm/versions/node/v22.22.0/bin/npm run start:server
Restart=on-failure

[Install]
WantedBy=multi-user.target
