# deploy/ihostmc-backend.service
[Unit]
Description=iHostMC Backend API
After=network.target postgresql.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/iHostMC/backend
Environment=NODE_ENV=production
Environment=PORT=3010
EnvironmentFile=/opt/iHostMC/backend/.env
# Use full path to node if needed (e.g. nvm: /home/ubuntu/.nvm/versions/node/v20.x/bin/node)
ExecStart=/home/ubuntu/.nvm/versions/node/v22.22.0/bin/node dist/index.js
Restart=on-failure

[Install]
WantedBy=multi-user.target
